# net-proxy/shadowsocks-libev/files/shadowsocks-initd
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

SS_PID="--pidfile /var/run/shadowsocks.pid"
SS_CONF="/etc/shadowsocks.json"

depend() {
    need net
}

checkconfig() {
    if [ ! -f ${SS_CONF} ]; then
        ewarn "${SS_CONF} does not exist."
    fi

    return 0
}

start() {
    checkconfig || return 1

    ebegin "Starting Shadowsocks"
    start-stop-daemon --start --quiet \
    $SS_PID --make-pid \
    --exec /usr/bin/ss-local -- -c /etc/shadowsocks.json &
    eend $?
}

stop() {
    ebegin "Stopping Shadowsocks"
    start-stop-daemon --stop --quiet ${SS_PID}
    eend $?
}

